# voice-ptt/com.sky.voice-ptt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sky.voice-ptt</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/sky/Work/Configs/voice-ptt/launch_voice_ptt.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/sky/Work/Configs/voice-ptt</string>

  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
  </array>

  <key>ProcessType</key>
  <string>Interactive</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/voice-ptt.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/voice-ptt.log</string>
</dict>
</plist>
