# megasync/mega-cmd-server.service
[Unit]
Description=MEGAcmd background server
Documentation=https://github.com/meganz/MEGAcmd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/mega-cmd-server --do-not-log-to-stdout
ExecStop=/usr/bin/mega-quit
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
